# files/etc/rc.local
#!/bin/sh -e
/etc/init.d/modemcontrol start
exit 0
