# files/etc/hotplug.d/iface/99-modem
[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "modem" ] && {
    logger -t modem "Modem interface is up"
}
