# files/etc/hotplug.d/usb/20-modem
#!/bin/sh
# Deteksi saat modem USB dicolok atau dicabut
[ "$ACTION" = "add" ] && logger -t modem "Modem plugged in: $PRODUCT"
[ "$ACTION" = "remove" ] && logger -t modem "Modem removed: $PRODUCT"
